# templates/haproxy/sync_certificates.cron
* * * * * root flock -n /var/lock/certsync.lock /etc/haproxy/sync_certificates.sh "<%= @certificate_sync_source %>"
